# oldtree/current/daemons/gpm/filelist
.FILELIST
.PKGINFO
etc/
etc/conf.d/
etc/conf.d/gpm
etc/rc.d/
etc/rc.d/gpm
usr/
usr/bin/
usr/bin/disable-paste
usr/bin/gpm-root
usr/bin/hltest
usr/bin/mev
usr/bin/mouse-test
usr/include/
usr/include/gpm.h
usr/lib/
usr/lib/libgpm.a
usr/lib/libgpm.so
usr/lib/libgpm.so.1.19.0
usr/man/
usr/man/man1/
usr/man/man1/gpm-root.1.gz
usr/man/man1/mev.1.gz
usr/man/man1/mouse-test.1.gz
usr/man/man7/
usr/man/man7/gpm-types.7.gz
usr/man/man8/
usr/man/man8/gpm.8.gz
usr/sbin/
usr/sbin/gpm
